# salt/nova-xvpvncproxy/etc/supervisor/conf.d/nova-xvpvncproxy.conf
[program:nova-xvpvncproxy]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-xvpvncproxy -- --config-file=/etc/nova/nova.conf
