# salt/mongodb-server/etc/supervisor/conf.d/mongodb-server.conf
[program:mongodb-server]
autorestart =  true
autostart   =  true
command     =  mongod -f /etc/mongodb.conf
