# salt/nova-api-metadata/etc/supervisor/conf.d/nova-api-metadata.conf
[program:nova-api-metadata]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-api-metadata -- --config-file=/etc/nova/nova.conf
