# salt/swift-proxy/etc/supervisor/conf.d/swift-proxy.conf
[program:swift-proxy]
autorestart =  true
autostart   =  true
command     = /usr/bin/swift-init proxy-server start
