# salt/nova-consoleauth/etc/supervisor/conf.d/nova-consoleauth.conf
[program:nova-consoleauth]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-consoleauth -- --config-file=/etc/nova/nova.conf
