# salt/nova-api-ec2/etc/supervisor/conf.d/nova-api-ec2.conf
[program:nova-api-ec2]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-api-ec2 -- --config-file=/etc/nova/nova.conf
