# salt/salt-minion/etc/supervisor/conf.d/salt-minion.conf
[program:salt-minion]
autorestart =  true
autostart   =  true
command     =  salt-minion
