# salt/salt-master/etc/supervisor/conf.d/salt-master.conf
[program:salt-master]
autorestart =  true
autostart   =  true
command     =  salt-master
