# salt/nova-cells/etc/supervisor/conf.d/nova-cells.conf
[program:nova-cells]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-cells -- --config-file=/etc/nova/nova.conf
