# salt/ceilometer-agent-notification/etc/supervisor/conf.d/ceilometer-agent-notification.conf
[program:ceilometer-agent-notification]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid ceilometer --exec /usr/bin/ceilometer-agent-notification
