# salt/cloudera-cm5-server/etc/systemd/system/cloudera-cm5-server.service
[Unit]
Description=cloudera-cm5-server
After=docker.service
Requires=docker.service

[Service]
ExecStart=/bin/bash -c 'docker start cloudera-cm5-server || docker run --detach=true --link=salt-master:salt --name=cloudera-cm5-server cloudera-cm5-server'
ExecStop=/bin/bash  -c 'docker stop  cloudera-cm5-server'

[Install]
WantedBy=multi-user.target
