# salt/glance-api/etc/supervisor/conf.d/glance-api.conf
[program:glance-api]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid glance --chdir /var/lib/glance --name glance-api --exec /usr/bin/glance-api
