# salt/nova-cert/etc/supervisor/conf.d/nova-cert.conf
[program:nova-cert]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-cert -- --config-file=/etc/nova/nova.conf
