# salt/skydns1/etc/supervisor/conf.d/skydns1.conf
[program:skydns1]
autorestart =  true
autostart   =  true
command     =  skydns -data=/var/lib/skydns1/data
