# salt/nova-objectstore/etc/supervisor/conf.d/nova-objectstore.conf
[program:nova-objectstore]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-objectstore -- --config-file=/etc/nova/nova.conf
