# salt/nova-novncproxy/etc/supervisor/conf.d/nova-novncproxy.conf
[program:nova-novncproxy]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-novncproxy -- --config-file=/etc/nova/nova.conf
