# salt/salt-syndic/etc/supervisor/conf.d/salt-syndic.conf
[program:salt-syndic]
autorestart =  true
autostart   =  true
command     =  salt-syndic
