# salt/swift-account/etc/supervisor/conf.d/swift-account.conf
[program:swift-account]
autorestart =  true
autostart   =  true
command     = /usr/bin/swift-init account-server start
