# salt/fleet/etc/supervisor/conf.d/fleet.conf
{% set ipv4 = salt['config.get']('fqdn_ip4') -%}
[program:fleet]
autorestart =  true
autostart   =  true
command     =  fleet
