# salt/nova-api/etc/supervisor/conf.d/nova-api.conf
[program:nova-api]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-api -- --config-file=/etc/nova/nova.conf
