# salt/heat-api-cloudwatch/etc/supervisor/conf.d/heat-api-cloudwatch.conf
[program:heat-api-cloudwatch]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid heat --exec /usr/bin/heat-api-cloudwatch
