# salt/swift-container/etc/supervisor/conf.d/swift-container.conf
[program:swift-container]
autorestart =  true
autostart   =  true
command     = /usr/bin/swift-init container-server start
