# salt/ironic-api/etc/supervisor/conf.d/ironic-api.conf
[program:ironic-api]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid ironic --exec /usr/bin/ironic-api -- --config-file=/etc/ironic/ironic.conf
