# salt/ceph-mon/etc/supervisor/conf.d/ceph-mon.conf
{% set cluster = salt['grains.get']('environment', 'ceph') -%}
{% set host    = salt['config.get']('host') -%}
[program:ceph-mon]
autorestart =  true
autostart   =  true
command     =  ceph-mon --cluster={{ cluster }} --id={{ host }}
