# salt/skydns2/etc/supervisor/conf.d/skydns2.conf
[program:skydns2]
autorestart =  true
autostart   =  true
command     =  bash --login -c skydns2
