# salt/etcd/etc/supervisor/conf.d/etcd.conf
[program:etcd]
autorestart =  true
autostart   =  true
command     =  etcd
