# salt/coreos-ipxe-server/etc/supervisor/conf.d/coreos-ipxe-server.conf
[program:coreos-ipxe-server]
autorestart =  true
autostart   =  true
command     = /opt/coreos-ipxe-server/bin/coreos-ipxe-server
