# salt/ceilometer-agent-central/etc/supervisor/conf.d/ceilometer-agent-central.conf
[program:ceilometer-agent-central]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid ceilometer --exec /usr/bin/ceilometer-agent-central
