# salt/nova-spiceproxy/etc/supervisor/conf.d/nova-spiceproxy.conf
[program:nova-spiceproxy]
autorestart =  true
autostart   =  true
command     =  start-stop-daemon --start --chuid nova --exec /usr/bin/nova-spicehtml5proxy -- --config-file=/etc/nova/nova.conf
