# salt/graphite-web/etc/supervisor/conf.d/graphite-web.conf
{% set graphite_py  = salt['config.get']('saltpath') + '/../graphite' -%}
{% set graphite_grp = salt["config.get"]("graphite-web:group:name") -%}
{% set graphite_usr = salt["config.get"]("graphite-web:user:name") -%}
[program:graphite-web]
autostart   =  true
autorestart =  true
command     =  gunicorn_django -u {{ graphite_usr }} -g {{ graphite_grp }} -b unix:/var/run/gunicorn/graphite-web.sock --log-file=/var/log/gunicorn/graphite-web.log {{ graphite_py }}/settings.py
user        =  {{ graphite_usr }}
